# runit/log/run
#!/bin/sh
mkdir /tmp/pan-bot
exec svlogd -tt /tmp/pan-bot
